# config/autostart/syncthing.desktop
[Desktop Entry]
Exec=syncthing --no-browser
Icon=
Name=syncthing
Path=
Terminal=False
Type=Application
